# config/systemd/local-orchestrator-worker.service.example
[Unit]
Description=Muel Local Orchestrator Worker
After=network.target

[Service]
Type=simple
User=muel
WorkingDirectory=/opt/muel/discord-news-bot
EnvironmentFile=/opt/muel/discord-news-bot/config/env/local-orchestrator-worker.env
ExecStart=/usr/bin/node --import tsx /opt/muel/discord-news-bot/scripts/agent-role-worker.ts --role=local-orchestrator
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/muel/discord-news-bot

[Install]
WantedBy=multi-user.target