# config/systemd/obsidian-lore-sync.service
[Unit]
Description=Obsidian Lore Sync to Supabase (periodic)
After=network.target unified-mcp-http.service
Requires=unified-mcp-http.service

[Service]
Type=oneshot
User=muel
WorkingDirectory=/opt/muel/muel-platform
Environment=HOME=/opt/muel
EnvironmentFile=-/opt/muel/muel-platform/.env
ExecStart=/usr/bin/node --import tsx /opt/muel/muel-platform/scripts/sync-obsidian-lore.ts

[Install]
WantedBy=multi-user.target
