# config/systemd/opencode-local-worker.service.example
[Unit]
Description=Muel Opencode Local Worker
After=network.target

[Service]
Type=simple
User=muel
WorkingDirectory=/opt/muel/muel-platform
EnvironmentFile=/opt/muel/muel-platform/config/env/opencode-worker.gcp.env
ExecStart=/usr/bin/node /opt/muel/muel-platform/scripts/opencode-local-worker.mjs
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/muel/muel-platform

[Install]
WantedBy=multi-user.target
