# config/systemd/unified-mcp-http.service
[Unit]
Description=MCP Unified HTTP Server (shared full-catalog MCP)
After=network.target obsidian-headless.service
Wants=obsidian-headless.service

[Service]
Type=simple
User=muel
WorkingDirectory=/opt/muel/shared-mcp-runtime
Environment=HOME=/opt/muel
Environment=NODE_ENV=production
Environment=MCP_HTTP_PORT=8850
EnvironmentFile=-/opt/muel/shared-mcp-runtime/.env
EnvironmentFile=-/opt/muel/shared-mcp-runtime/config/env/unified-mcp.gcp.env
ExecStart=/bin/bash -lc 'set -a; if [ -f /opt/muel/shared-mcp-runtime/.env ]; then . /opt/muel/shared-mcp-runtime/.env; fi; if [ -f /opt/muel/shared-mcp-runtime/config/env/unified-mcp.gcp.env ]; then . /opt/muel/shared-mcp-runtime/config/env/unified-mcp.gcp.env; fi; set +a; exec /usr/bin/node --import tsx /opt/muel/shared-mcp-runtime/scripts/unified-mcp-http.ts'
Restart=on-failure
RestartSec=10
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
